# home/config/systemd/user/sxmo_modemmonitor.service
[Unit]
Description=Notifications service
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/sxmo_modemmonitor.sh

[Install]
WantedBy=sway-session.target
